# files/rc.d/rc.local
#!/bin/bash

# Run all cloudconfig fixes every time. 
for i in /etc/cloudconfig/fixes.d/*.sh; do
    if [ -r "$i" ]; then
        . "$i" >> /var/log/cloudconfig.log
    fi
done

touch /var/lock/subsys/local

