# files/rc.d/rc.local
#!/bin/bash
#
# Give this host a random hostname if it doesn't get something
# unique from DHCP. 
#

SN=`hostname --short`
if [ "${SN}x" == "localhostx" ] ; then
   hostname $RANDOM.localdomain
fi


if [ ! -d /root/.ssh ]; then
  mkdir -p /root/.ssh
  cp /etc/ssh/authorized_keys /root/.ssh/
  chmod go-rwx /root/.ssh
fi

touch /var/lock/subsys/local

