# files/profile.d/osg.sh
#
# Handle setup for OSG worker node.  
#
export OSG_GRID=/etc/osg/wn-client
export OSG_APP=/home/osg/app
export OSG_DATA=/home/osg/data
export OSG_WN_TMP=/home/osg/data
export LCG_GFAL_INFOSYS="is.grid.iu.edu:2170"

